# private/zuup-os/security/kiosk/zuup-kiosk.service
# ZUUP-OS kiosk unit (spec §7.4, §7.5) — systemd → Cage → Firefox, no escape.
#
# Cage is a single-surface Wayland compositor: it refuses every additional
# wl_surface, so no second window, overlay, or screenshot-of-another-window can
# exist. Firefox launches --kiosk at the Edge. The hardening directives below
# are the systemd half of the §7.5 runtime enforcement (the AppArmor + seccomp
# halves live alongside this file). Restart=always: there is no way to "exit".

[Unit]
Description=ZUUP-OS locked examination surface
# The Gate renders only after the boot is attested (§7.1) and the biometric
# daemon is serving — and never without the tunnel (INV-3, INV-10).
# Cage needs a seat manager (seatd) to take the DRM display + input — logind is
# masked on this image, so seatd is the one that hands cage the seat.
After=zuup-network.target zuup-attest.service zuup-biometric.service seatd.service
Requires=zuup-network.target zuup-attest.service seatd.service

[Service]
Type=simple
# Cage owns the only display surface; Firefox is its only child. The launcher
# picks the URL from this terminal's Edge-reported capability (§7.6) so the
# ADMIN_STATION opens the Centre Admin portal, an INVIGILATOR_STATION/
# CANDIDATE_SEAT opens the exam-terminal Gate, and anything else fails closed.
# Cage needs an XDG_RUNTIME_DIR seat and Firefox a writable HOME — both RAM-only
# (INV-2): /run/firefox via RuntimeDirectory, the profile under PrivateTmp's /tmp.
Environment=XDG_RUNTIME_DIR=/run/firefox
Environment=HOME=/tmp/firefox
ExecStartPre=/bin/mkdir -p /tmp/firefox
ExecStart=/usr/lib/zuup/zuup-kiosk-launch.sh
Restart=always
RestartSec=1

# ── §7.5 runtime hardening ───────────────────────────────────────────────────
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ProtectKernelModules=yes
ProtectKernelTunables=yes
ProtectControlGroups=yes
ProtectClock=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
RestrictAddressFamilies=AF_UNIX AF_INET
LockPersonality=yes
MemoryDenyWriteExecute=yes
SystemCallFilter=@system-service
SystemCallArchitectures=native
SystemCallErrorNumber=EPERM
# seccomp allow-list (the explicit kill-list complement, §7.5)
SystemCallFilter=~@debug @mount @reboot @swap @obsolete
# Writable surfaces are tmpfs only (RAM; INV-2). RuntimeDirectory creates
# /run/firefox at start (0700, auto-removed on stop); PrivateTmp gives a fresh
# private /tmp tmpfs for the Firefox profile. These EXIST before the sandbox is
# wired, unlike a ReadWritePaths bind of a not-yet-created /tmp/firefox — which
# is what made systemd fail at 226/NAMESPACE and crash-loop the kiosk.
RuntimeDirectory=firefox
RuntimeDirectoryMode=0700
PrivateTmp=yes
# VT switching is disabled at logind; Cage takes the only seat.
TTYVTDisallocate=yes

[Install]
WantedBy=zuup-session.target
